# app/deploy/dev.apifor.executor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd (macOS) — registra o executor como serviço de fundo do usuário.
     Instale em ~/Library/LaunchAgents/ e: launchctl load -w <plist> -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.apifor.executor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/apifor-executor</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>APIFOR_HOME</key>
        <string>/usr/local/var/apifor</string>
        <key>APIFOR_WORKDIR</key>
        <string>/usr/local/var/apifor/work</string>
        <key>CEREBRO_HTTP</key>
        <string>https://api.apifor.dev</string>
        <key>CEREBRO_GRPC</key>
        <string>https://api.apifor.dev:9090</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/apifor-executor.log</string>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/apifor-executor.err</string>
</dict>
</plist>
